# dockerfiles/graphite-stack/supervisor.d/api.ini
[program:carbonzipper]
command=carbonzipper -c /openmetric/conf/zipper.conf -logdir /openmetric/log/
priority=500
autorestart=true

[program:carbonapi]
command=carbonapi-wrapper /openmetric/conf/api.conf
priority=999
autorestart=true
