# dockerfiles/graphite-stack/supervisor.d/carbon.ini
[program:carbon]
command=go-carbon --config /openmetric/conf/carbon.conf
priority=1
autorestart=true
