# dockerfiles/openmetric/supervisor.d/relay.ini
[program:relay]
command=carbon-c-relay -f /openmetric/conf/relay.conf -l /openmetric/log/relay.log -c '-_:#@'
autorestart=true
